# serverstats/server-stats.service
# Oneshot service that takes one server-statistics sample. The companion timer fires it
# on the configured interval. It runs on the host as the rootless-podman user (not as a
# container) because only the host can read the pod's cgroup counters, disk IOPS and
# network egress. The script reaches the database with `podman exec postgresql`, so it
# needs nothing published on a port.
#
# Installed under ~/.config/systemd/user/ by install.sh, alongside the quadlet-generated
# units, so `systemctl --user` and `journalctl --user` manage and show it like the rest.
# ${APP_NAME} is rendered by the installer (envsubst), the same way the quadlets are.
[Unit]
Description=${APP_NAME} server-statistics sample
# The sample reaches the database through the postgresql container, so only attempt it
# once the pod is up. This is a hint, not a hard dependency: a tick that runs while the
# database is down simply fails and the next one succeeds.
After=main-pod.service
Wants=main-pod.service

[Service]
Type=oneshot
# The collector derives its config and state paths from the app name.
Environment=APP_NAME=${APP_NAME}
ExecStart=%h/.config/${APP_NAME}/serverstats/collect.sh
