# quadlets/backup.container
# A one-shot database backup. Triggered by backup.timer, it joins the pod to
# reach the database on localhost, dumps every database and role with pg_dumpall to the
# backup volume, then exits. Type=oneshot so the timer treats a clean exit as success.
#
# This is deliberately decoupled from deployment: podman auto-update can roll back a bad
# image, but never the database, so a regular dump is the only safeguard against a
# migration that corrupts data. Restore a dump with:
#   podman run --rm -i --pod gefieder -v gefieder-backup:/backup \
#     ${REGISTRY}/postgresql:${IMAGE_TAG} \
#     sh -c 'psql -h localhost -U ${SUPERUSER_NAME} -d postgres < /backup/<file>.sql'
[Unit]
Description=Gefieder database backup (pg_dumpall)
After=postgresql.service
Requires=postgresql.service

[Container]
ContainerName=${APP_NAME}-backup
Pod=${APP_NAME}.pod
Image=${REGISTRY}/postgresql:${IMAGE_TAG}
Volume=${APP_NAME}-backup.volume:/backup
Secret=superuser_password
# PGPASSWORD lets pg_dumpall authenticate non-interactively over the pod's localhost.
Exec=sh -c 'PGPASSWORD=$(cat /run/secrets/superuser_password) pg_dumpall -h localhost -U ${SUPERUSER_NAME} --clean --if-exists > /backup/db-$(date +%%Y%%m%%d-%%H%%M%%S).sql'

[Service]
Type=oneshot
