# quadlets/grafana.container
# The Grafana dashboards, with the database pre-configured as a read-only data source.
# The provisioning is baked into the custom image install.sh loaded from the release.
[Unit]
Description=Grafana dashboards
After=postgresql.service

[Container]
ContainerName=grafana
Pod=main.pod
Image=${REGISTRY}/grafana:${IMAGE_TAG}
# APP_NAME names the provisioned data source. The admin user logs in with the superuser
# name and password, the same credentials as the PostgreSQL and Django superuser.
Environment=APP_NAME=${APP_NAME}
Environment=GF_SECURITY_ADMIN_USER=${SUPERUSER_NAME}
Environment=GF_SECURITY_ADMIN_PASSWORD__FILE=/run/secrets/superuser_password
# Serve Grafana from the /GRAFANA_PATH subpath behind the proxy so its internal links
# and assets resolve correctly. This stays here rather than in custom.ini because it
# depends on the build-time GRAFANA_PATH variable.
Environment=GF_SERVER_ROOT_URL=%(protocol)s://%(domain)s/${GRAFANA_PATH}/
Environment=GF_SERVER_SERVE_FROM_SUB_PATH=true
# Persist the log into the data volume (in addition to the container's stdout), so a
# crash leaves a log on disk. GF_PATHS_LOGS points at a subdir of the mounted volume;
# the grafana user writes it and the rootless container maps it back to the podman user.
# These stay here as GF_* env, not in custom.ini: the image's entrypoint passes paths.logs
# and log.mode on grafana-server's command line, which overrides the ini file but honours a
# pre-set GF_* environment variable.
# The value is quoted because systemd's Environment= otherwise splits on the space and
# would set GF_LOG_MODE=console only, dropping the "file" sink.
Environment="GF_LOG_MODE=console file"
Environment=GF_PATHS_LOGS=/var/lib/grafana/log
Volume=grafana_data.volume:/var/lib/grafana
# superuser_password for the admin login, grafana_password for the data source.
Secret=superuser_password
Secret=grafana_password
HealthCmd=wget -qO /dev/null http://localhost:3000/${GRAFANA_PATH}/api/health
HealthInterval=10s
HealthTimeout=5s
HealthRetries=10
HealthStartPeriod=20s

[Service]
Restart=always

[Install]
WantedBy=default.target
