# quadlets/postgresql.container
# The PostgreSQL database holding the engineering, analytics and application data.
#
# The image is the one install.sh loaded from the release (REGISTRY/IMAGE_TAG names it).
# Notify=healthy ties "service started" to the healthcheck below, so units ordered after
# this one wait for the database to actually accept connections.
[Unit]
Description=Gefieder PostgreSQL database

[Container]
ContainerName=postgresql
Pod=main.pod
Image=${REGISTRY}/postgresql:${IMAGE_TAG}
Notify=healthy
Environment=POSTGRES_DB=postgres
Environment=POSTGRES_USER=${SUPERUSER_NAME}
Environment=POSTGRES_PASSWORD_FILE=/run/secrets/superuser_password
Volume=postgresql_data.volume:/var/lib/postgresql/data
Secret=superuser_password
Secret=crudman_password
Secret=sqlmesh_password
Secret=grafana_password
HealthCmd=pg_isready -U ${SUPERUSER_NAME} -d postgres
HealthInterval=5s
HealthTimeout=5s
HealthRetries=10
HealthStartPeriod=10s

[Service]
Restart=always

[Install]
WantedBy=default.target
